# templates/etc/systemd/system/patroni.service.d/override.conf.j2
[Service]
# Resets ExecStart before override its value
ExecStart=
ExecStart={{ patroni_bin_dir }}/patroni {{ patroni_config_dir }}/
